# sync.cron
SHELL=/bin/bash
* * * * * tivvit /home/tivvit/git/stalker/sync.sh pc
